# qr3_87.conf
[program:qr3_87]
command=python3 /home/curator3/qr3_87/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator3/qr3_87/test.err.log
stdout_logfile=/home/curator3/qr3_87/test.err.log
